# scripts/selfheal/systemd/atp-health-alert.service
[Unit]
Description=ATP health snapshot Telegram alert (check log, send if rule triggers)
After=network-online.target atp-health-snapshot.service

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/crypto-2.0
ExecStart=/home/ubuntu/crypto-2.0/scripts/diag/health_snapshot_telegram_alert.sh
User=root

[Install]
WantedBy=multi-user.target
