# scripts/selfheal/systemd/atp-health-snapshot.service
[Unit]
Description=ATP hourly health snapshot to /var/log/atp/health_snapshots.log
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/crypto-2.0
ExecStart=/home/ubuntu/crypto-2.0/scripts/diag/health_snapshot_log.sh
# Root so /var/log/atp can be created and written
User=root

[Install]
WantedBy=multi-user.target
