# scripts/health_monitor.service
[Unit]
Description=Automated Trading Platform Health Monitor
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/crypto-2.0
ExecStart=/home/ubuntu/crypto-2.0/scripts/health_monitor.sh
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/crypto-2.0/logs/health_monitor.log
StandardError=append:/home/ubuntu/crypto-2.0/logs/health_monitor.error.log

[Install]
WantedBy=multi-user.target

