# scripts/automation/systemd/jarvis-openclaw-guard.timer
[Unit]
Description=Run OpenClaw regression guard every 6 hours

[Timer]
OnBootSec=300
OnCalendar=*-*-* 00,06,12,18:00:00 UTC
AccuracySec=5m
Persistent=true

[Install]
WantedBy=timers.target
