# scripts/automation/systemd/jarvis-daily-report.service
[Unit]
Description=Jarvis daily production report (read-only)
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=@REPO_ROOT@
EnvironmentFile=-@REPO_ROOT@/secrets/runtime.env
EnvironmentFile=-@REPO_ROOT@/.env.aws
ExecStart=@PYTHON@ @REPO_ROOT@/scripts/automation/daily_report.py
User=root

[Install]
WantedBy=multi-user.target
