# scripts/automation/systemd/jarvis-health-check.service
[Unit]
Description=Jarvis health check automation (read-only)
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=@REPO_ROOT@
EnvironmentFile=-@REPO_ROOT@/secrets/runtime.env
EnvironmentFile=-@REPO_ROOT@/.env.aws
ExecStart=@PYTHON@ @REPO_ROOT@/scripts/automation/health_check.py
User=root

[Install]
WantedBy=multi-user.target
